# deploy/container/s6-rc.d/init-host-keys/up
/bin/sh /etc/s6-overlay/scripts/init-host-keys.sh
