# deploy/container/s6-rc.d/dashboard/finish
#!/bin/sh
# Called by s6 when the dashboard's run script exits, for any reason.
# See sshd/finish for the args contract.
exit_code="$1"
signal="$2"
if [ "$exit_code" = "256" ]; then
  reason="killed by signal $signal"
else
  reason="exited with code $exit_code"
fi
echo "[$(date -Iseconds)] nbg-supervise: service dashboard $reason - s6 will restart" >&2
exit 0
